# scripts/maintenance/portfolio_cleanup.timer
# ==============================================================================
# SPOKE - PORTFOLIO CLEANUP TIMER
# ==============================================================================
# Description: systemd user timer — runs portfolio_cleanup.service daily
# Author: Matt Barham
# Created: 2026-03-08
# Modified: 2026-04-22
# Version: 1.0.1
# Host: Your Server
# ==============================================================================
# Type: systemd unit (timer)
# Component: Spoke / maintenance
# ==============================================================================

[Unit]
Description=Portfolio Database Cleanup Timer

[Timer]
# Run daily at 3:00 AM
OnCalendar=*-*-* 03:00:00
Persistent=true
RandomizedDelaySec=300

[Install]
WantedBy=timers.target
