# scripts/maintenance/spoke_log_analysis.service
[Unit]
Description=Spoke Daily Log Analysis Report

[Service]
Type=oneshot
# TODO: Edit ExecStart path for your installation
# ExecStartPre polls Docker readiness (user units cannot depend on system docker.service)
ExecStartPre=/bin/bash -c 'for i in $(seq 1 30); do docker info &>/dev/null && exit 0; sleep 2; done; exit 1'
ExecStart=/path/to/spoke/scripts/maintenance/spoke_log_analysis.sh
StandardOutput=journal
StandardError=journal
SyslogIdentifier=spoke-log-analysis
TimeoutStartSec=600

# Security settings
PrivateTmp=yes
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/tmp

[Install]
WantedBy=default.target
