# scripts/maintenance/spoke-backup.timer
# ==============================================================================
# SPOKE - CONFIG BACKUP TIMER
# ==============================================================================
# Description: systemd user timer — runs spoke-backup.service daily
# Author: Matt Barham
# Created: 2026-04-06
# Modified: 2026-04-22
# Version: 1.0.1
# Host: Your Server
# ==============================================================================
# Type: systemd unit (timer)
# Component: Spoke / maintenance (backup)
# ==============================================================================

[Unit]
Description=Spoke Config Backup Timer
Requires=spoke-backup.service

[Timer]
# Run daily at 2:00 AM local time
OnCalendar=*-*-* 02:00:00
Persistent=true
RandomizedDelaySec=300

[Install]
WantedBy=timers.target
