# scripts/maintenance/spoke-backup.timer
[Unit]
Description=Spoke Config Backup Timer
Requires=spoke-backup.service

[Timer]
# Run daily at 2:00 AM local time
OnCalendar=*-*-* 02:00:00
Persistent=true
RandomizedDelaySec=300

[Install]
WantedBy=timers.target
