# scripts/maintenance/update_duckdns.service
# ==============================================================================
# SPOKE - DUCKDNS UPDATER SERVICE
# ==============================================================================
# Description: systemd service — runs update_duckdns.sh to refresh DDNS record
# Author: Matt Barham
# Created: 2026-03-23
# Modified: 2026-04-22
# Version: 1.0.1
# Host: Your Server
# ==============================================================================
# Type: systemd unit (service)
# Component: Spoke / maintenance (DDNS)
# ==============================================================================

[Unit]
Description=Update DuckDNS DDNS record with current public IP
Requires=network-online.target
After=network-online.target

[Service]
Type=oneshot
# TODO: Edit User, Group, and ExecStart path for your installation
User=your-user
Group=docker
ExecStart=/path/to/spoke/scripts/maintenance/update_duckdns.sh
StandardOutput=journal
StandardError=journal
SyslogIdentifier=duckdns-updater

# Security settings
PrivateTmp=yes
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/tmp

[Install]
WantedBy=multi-user.target
