# scripts/maintenance/update_duckdns.timer
# ==============================================================================
# SPOKE - DUCKDNS UPDATER TIMER
# ==============================================================================
# Description: systemd timer — triggers update_duckdns.service every 10 min
# Author: Matt Barham
# Created: 2026-03-16
# Modified: 2026-04-22
# Version: 1.0.1
# Host: Your Server
# ==============================================================================
# Type: systemd unit (timer)
# Component: Spoke / maintenance (DDNS)
# ==============================================================================

[Unit]
Description=Update DuckDNS DDNS record every 10 minutes

[Timer]
OnBootSec=60
OnUnitActiveSec=600
Persistent=true

[Install]
WantedBy=timers.target
